# titties-on-boot.service
[Unit]
Description=Running the TITTIES Telegram bot and speaker engine.
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
WorkingDirectory=/home/pi/t-i-t-t-i-e-s/standalone
ExecStart=/home/pi/.deno/bin/deno task run
Restart=always
RestartSec=5
User=pi

[Install]
WantedBy=multi-user.target
